# deploy/schulbefragung.service
[Unit]
Description=Schulbefragung Web-App
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/schulbefragung
EnvironmentFile=/opt/schulbefragung/.env
ExecStart=/opt/schulbefragung/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 1
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
